# private_dot_config/systemd/user/opencode.service
[Unit]
Description=OpenCode server
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/opencode-server.fish
Restart=on-failure
RestartSec=2
Environment=OP_ACCOUNT=createxyz.1password.com

[Install]
WantedBy=default.target
